# deployment/scrapyd.conf
[program:scrapyd]
directory=/root/pension_crawler
command=/usr/local/bin/scrapyd
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=syslog
stderr_logfile=syslog
